# system/etc/init/fat.conf
start on started udev 

script
sleep 2
mount -o remount,rw /
mount |tee /tmp/mount.log
if [ ! -e /opt/mfg/ui/log/fat.end ]; then
    cd /opt/mfg/ui
    ./fat.sh > /dev/tty1 < /dev/tty1
fi
end script
